# deploy/systemd/spotfam-deploy.service
[Unit]
Description=SpotFamServ Pull-Deploy (tag-getriggert)
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=lars
WorkingDirectory=/home/lars/SpotFamServ
# Public Health-Endpunkt: greift auch, falls eine aeltere in-memory Skriptversion noch
# den alten Default (/api/v1/profiles, jetzt auth-geschuetzt -> 401) verwenden wuerde.
Environment=HEALTH_URL=http://localhost:8080/api/v1/health
ExecStart=/home/lars/SpotFamServ/deploy/pi-deploy.sh
# Logs landen im Journal: journalctl -u spotfam-deploy.service
